# config/crontab
* * * * * /usr/bin/python3 /root/get_curr_weather.py 2>&1